# pkg-build/scripts/postinstall
#!/bin/bash
# postinstall script for Screenshot Renamer
# Runs as root after pkg payload is installed to /usr/local/share/screenshot-renamer

set -e

# Detect the console user (not root)
CURRENT_USER=$( /usr/bin/stat -f "%Su" /dev/console )
CURRENT_USER_HOME=$( /usr/bin/dscl . -read /Users/"$CURRENT_USER" NFSHomeDirectory | /usr/bin/awk '{print $NF}' )
CURRENT_USER_UID=$( /usr/bin/id -u "$CURRENT_USER" )

STAGING="/usr/local/share/screenshot-renamer"
SCREENSHOTS_DIR="$CURRENT_USER_HOME/Screenshots"
WATCHER_SCRIPT="$CURRENT_USER_HOME/Library/Scripts/rename-new-screenshots.sh"
LAUNCH_AGENT_PLIST="$CURRENT_USER_HOME/Library/LaunchAgents/nz.glen.screenshot-renamer.plist"
LAUNCH_AGENT_LABEL="nz.glen.screenshot-renamer"

# 1. Create Screenshots directory
mkdir -p "$SCREENSHOTS_DIR"
chown "$CURRENT_USER" "$SCREENSHOTS_DIR"

# 2. Set macOS to save screenshots there
/usr/bin/sudo -u "$CURRENT_USER" defaults write com.apple.screencapture location "$SCREENSHOTS_DIR"
killall SystemUIServer 2>/dev/null || true

# 3. Install the Shortcut (opens GUI prompt for user)
SHORTCUT_PATH="$STAGING/Rename Screenshot.shortcut"
if [ -f "$SHORTCUT_PATH" ]; then
    /usr/bin/sudo -u "$CURRENT_USER" open "$SHORTCUT_PATH"
fi

# 4. Install the watcher script
mkdir -p "$CURRENT_USER_HOME/Library/Scripts"
cat > "$WATCHER_SCRIPT" <<'WATCHER'
#!/bin/zsh
# Called by launchd WatchPaths whenever ~/Screenshots changes.
# Finds screenshots that still have the default macOS name and renames them
# via the "Rename Screenshot" Shortcut.

SCREENSHOTS_DIR="$HOME/Screenshots"
TIMESTAMP_FILE="$HOME/Library/Scripts/.rename-screenshot-lastrun"

if [[ ! -f "$TIMESTAMP_FILE" ]]; then
    touch -t 197001010000 "$TIMESTAMP_FILE"
fi

# Small sleep to let the screenshot finish writing to disk
sleep 2

find "$SCREENSHOTS_DIR" -maxdepth 1 -name "Screenshot *.png" -newer "$TIMESTAMP_FILE" | while IFS= read -r file; do
    shortcuts run "Rename Screenshot" --input-path "$file"
done

touch "$TIMESTAMP_FILE"
WATCHER
chmod +x "$WATCHER_SCRIPT"
chown "$CURRENT_USER" "$WATCHER_SCRIPT"

# 5. Install and load the Launch Agent
mkdir -p "$CURRENT_USER_HOME/Library/LaunchAgents"
cat > "$LAUNCH_AGENT_PLIST" <<PLIST
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>$LAUNCH_AGENT_LABEL</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>$WATCHER_SCRIPT</string>
    </array>
    <key>WatchPaths</key>
    <array>
        <string>$SCREENSHOTS_DIR</string>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>$CURRENT_USER_HOME/Library/Logs/rename-screenshot.log</string>
    <key>StandardErrorPath</key>
    <string>$CURRENT_USER_HOME/Library/Logs/rename-screenshot.log</string>
</dict>
</plist>
PLIST
chown "$CURRENT_USER" "$LAUNCH_AGENT_PLIST"

# Unload first if already running, then load as the user
/bin/launchctl asuser "$CURRENT_USER_UID" /bin/launchctl unload "$LAUNCH_AGENT_PLIST" 2>/dev/null || true
/bin/launchctl asuser "$CURRENT_USER_UID" /bin/launchctl load "$LAUNCH_AGENT_PLIST"

# 6. Clean up staging area
rm -rf "$STAGING"

exit 0
